# src/deb/control/postinst
#!/bin/bash
echo Route.put Post-Install Script Running!
chown -R routeput:routeput /usr/share/routeput
update-rc.d routeput defaults
update-rc.d routeput enable
/etc/init.d/routeput start
